# deploy/reachtolet-celery.service
[Unit]
Description=Reachtolet Celery worker (billboard tracking + notifications)
After=network.target redis-server.service
Wants=redis-server.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/Reachtolet_backend
Environment=PUBLIC_BASE_URL=http://16.16.160.64:8000
Environment=CELERY_BROKER_URL=redis://127.0.0.1:6379/0
Environment=DB_CONN_MAX_AGE=120
ExecStart=/home/ubuntu/Reachtolet_backend/.venv/bin/celery -A core worker --loglevel=info --concurrency=2
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
