# deploy/reachtolet-daphne.service
[Unit]
Description=Reachtolet Django ASGI (Daphne + Socket.IO)
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/Reachtolet_backend
EnvironmentFile=-/etc/reachtolet/reachtolet.env
ExecStart=/home/ubuntu/Reachtolet_backend/.venv/bin/daphne -b 0.0.0.0 -p 8000 --proxy-headers core.asgi:application \
  --access-log /home/ubuntu/Reachtolet_backend/logs/daphne-access.log \
  --verbosity 1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
